# server/systemd/bourmet-audit-archive.timer
[Unit]
Description=Archivia ogni mese i log amministrativi vecchi

[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=30m
Unit=bourmet-audit-archive.service

[Install]
WantedBy=timers.target
